# .config/systemd/user/wallpaper.service
[Unit]
Description=wallpaperscript

[Service]
Type=oneshot
Enviroment=$DISPLAY
ExecStart=/home/enzobarrett/wallpaper
 
